# crontab
* * * * * python /root/conv.py >/dev/null 2>&1
0 2 * * * rm -rf /audio/* >/dev/null 2>&1
5 2 * * * reboot >/dev/null 2>&1